# cix_proprietary/cix_proprietary-debs/cix-gpu-umd/DEBIAN/postinst
#!/bin/sh
set -e

# Enable the service to start on boot
systemctl enable load-modules.service || true

exit 0
